# pkgs/bigfred/server/supervisord/templates/supervisord.conf.tmpl
; Generated by BigFred — do not edit manually.
[unix_http_server]
file={{ .SocketPath }}
chmod=0700

[supervisord]
logfile={{ .LogDir }}/supervisord.log
pidfile={{ .PIDFile }}
childlogdir={{ .LogDir }}
nodaemon=false
user={{ .RunAsUser }}

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://{{ .SocketPath }}

{{- range .Groups }}

[group:{{ .Name }}]
programs={{ joinProgramNames .Programs }}
{{- range .Programs }}

[program:{{ .Name }}]
command=/bin/bash -c {{ shellQuote .Command }}
directory={{ $.ConfigDir }}
autostart={{ .Autostart }}
autorestart={{ .Autorestart }}
stdout_logfile={{ $.LogDir }}/{{ .Name }}.stdout.log
stderr_logfile={{ $.LogDir }}/{{ .Name }}.stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stopwaitsecs={{ orInt .StopWaitSecs 10 }}
startsecs={{ orInt .StartSecs 1 }}
{{- end }}
{{- end }}
